# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:node-api]
command=node /app/index.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="%(ENV_PORT)s"

[program:ml-service]
command=/opt/ml-venv/bin/uvicorn main:app --host 127.0.0.1 --port %(ENV_ML_SERVICE_PORT)s
directory=/app/ml-service
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=ML_SERVICE_API_KEY="%(ENV_ML_SERVICE_API_KEY)s"
